# server/rss/rssfeed_uwsgi.service
[Unit]
Description=Start uwsgi for hosting RSS feeds

[Service]
ExecStart=uwsgi --ini /caladrius/server/rss/wsgi_config.ini
User=http
Group=http

[Install]
WantedBy=multi-user.target
